# etc/systemd/system/containerstartupscript.service
[Unit]
Description=starts the containers up and updates them on startup
<<<<<<< HEAD

[Service]
ExecStart=/home/xilinx/StartupScript/ContainerScript.sh

[Install]
WantedBy=multi-user.target

=======
>>>>>>> 7faa424428c7b5ea7e1022f895f10449f2d7252c

[Service]
ExecStart=/home/xilinx/StartupScript/ContainerScript.sh

[Install]
WantedBy=multi-user.target
